# ops/systemd/sitbank-security-alerts.timer
[Unit]
Description=Run SITBank security alert checks every 5 minutes

[Timer]
OnActiveSec=5min
OnUnitActiveSec=5min
AccuracySec=30s
Unit=sitbank-security-alerts.service

[Install]
WantedBy=timers.target
